# deployment/systemd/redhouse-solar-prediction.timer
[Unit]
Description=RedHouse Solar Yield Prediction Timer
Requires=redhouse-solar-prediction.service

[Timer]
OnBootSec=3min
OnCalendar=hourly
OnCalendar=*:03:00
AccuracySec=1s

[Install]
WantedBy=timers.target
