# deployment/systemd/redhouse-checkwatt.timer
[Unit]
Description=RedHouse CheckWatt Data Collection Timer
Requires=redhouse-checkwatt.service

[Timer]
OnBootSec=3min
# Run every 5 minutes at :01, :06, :11, :16, :21, :26, :31, :36, :41, :46, :51, :56
OnCalendar=*:01,06,11,16,21,26,31,36,41,46,51,56:00
AccuracySec=1s

[Install]
WantedBy=timers.target
