# deployment/systemd/redhouse-aggregate-emeters-5min.timer
[Unit]
Description=RedHouse 5-Minute Energy Meter Aggregation Timer
Requires=redhouse-aggregate-emeters-5min.service
After=redhouse-checkwatt.timer

[Timer]
OnBootSec=5min
# IMPORTANT: If you change these times, also update deployment/smart_deploy.sh
# Run every 5 minutes at :03, :08, :13, :18, :23, :28, :33, :38, :43, :48, :53, :58
# Runs 2 minutes after CheckWatt (:01, :06, ...) to ensure all source data is collected.
# Shifted +1min from :02 to :03 to avoid overlap with wibatemp checkwatt_dataloader
# which runs at :01,:06,:11,... and does heavy InfluxDB reads for ~10s.
# TODO: Shift back to :02 once wibatemp is decommissioned.
OnCalendar=*:03,08,13,18,23,28,33,38,43,48,53,58:00
AccuracySec=1s

[Install]
WantedBy=timers.target
