# deployment/systemd/redhouse-backup.timer
[Unit]
Description=RedHouse Pi Backup Timer

[Timer]
OnCalendar=*-*-* 03:00:00
AccuracySec=60s
Persistent=true

[Install]
WantedBy=timers.target
