# deployment/systemd/redhouse-solar-prediction.service
[Unit]
Description=RedHouse Solar Yield Prediction
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/wibatemp
ExecStart=/usr/bin/python3 /home/pi/wibatemp/predict_solar_yield.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redhouse-solar-prediction

[Install]
WantedBy=multi-user.target
