# deployment/systemd/redhouse-evu-cycle.service
[Unit]
Description=RedHouse EVU Cycle (Prevent Direct Heating Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/wibatemp
ExecStart=/home/pi/wibatemp/mlp_cycle_evu_off_hourly.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redhouse-evu-cycle

[Install]
WantedBy=multi-user.target
