# deployment/systemd/redhouse-weather.service
[Unit]
Description=RedHouse Weather Forecast Collection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/opt/redhouse
ExecStart=/opt/redhouse/venv/bin/python /opt/redhouse/collect_weather.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redhouse-weather

[Install]
WantedBy=multi-user.target
