# deployment/systemd/redhouse-solar-prediction.service
[Unit]
Description=RedHouse Solar Yield Prediction
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/opt/redhouse
ExecStart=/opt/redhouse/venv/bin/python /opt/redhouse/predict_solar_yield.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redhouse-solar-prediction

[Install]
WantedBy=multi-user.target
