# deployment/systemd/redhouse-windpower.service
[Unit]
Description=RedHouse Wind Power Data Collection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/opt/redhouse
ExecStart=/opt/redhouse/venv/bin/python /opt/redhouse/collect_windpower.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redhouse-windpower

[Install]
WantedBy=multi-user.target
